# sigma/sysmon/process_creation/proc_creation_win_schtasks_folder_combos.yml
title: Schtasks From Suspicious Folders
ruletype: Sigma
author: Florian Roth (Nextron Systems)
date: 2022/04/15
description: Detects scheduled task creations that have suspicious action command
    and folder combinations
detection:
    SELECTION_1:
        EventID: 1
    SELECTION_2:
        Channel: Microsoft-Windows-Sysmon/Operational
    SELECTION_3:
        Image: '*\schtasks.exe'
    SELECTION_4:
        OriginalFileName: schtasks.exe
    SELECTION_5:
        CommandLine: '* /create *'
    SELECTION_6:
        CommandLine:
        - '*powershell*'
        - '*pwsh*'
        - '*cmd /c *'
        - '*cmd /k *'
        - '*cmd /r *'
        - '*cmd.exe /c *'
        - '*cmd.exe /k *'
        - '*cmd.exe /r *'
    SELECTION_7:
        CommandLine:
        - '*C:\ProgramData\\*'
        - '*%ProgramData%*'
    condition: ((SELECTION_1 and SELECTION_2) and (SELECTION_3 or SELECTION_4) and
        SELECTION_5 and SELECTION_6 and SELECTION_7)
falsepositives:
- Unknown
id: 8a8379b8-780b-4dbf-b1e9-31c8d112fefb
level: high
logsource:
    category: process_creation
    product: windows
modified: 2022/11/18
references:
- https://symantec-enterprise-blogs.security.com/blogs/threat-intelligence/lazarus-dream-job-chemical
status: experimental
tags:
- attack.execution
- attack.t1053.005

