# sigma/builtin/process_creation/proc_creation_win_schtasks_reg_loader.yml
title: Scheduled Task Executing Powershell Encoded Payload from Registry
ruletype: Sigma
author: '@Kostastsale, @TheDFIRReport, slightly modified by pH-T (Nextron Systems)'
date: 2022/02/12
description: Detects the creation of a schtask that executes a base64 encoded payload
    stored in the Windows Registry using PowerShell.
detection:
    SELECTION_1:
        EventID: 4688
    SELECTION_2:
        Channel: Security
    SELECTION_3:
        NewProcessName: '*\schtasks.exe'
    SELECTION_4:
        CommandLine: '*/Create*'
    SELECTION_5:
        CommandLine: '*/SC*'
    SELECTION_6:
        CommandLine: '*FromBase64String*'
    SELECTION_7:
        CommandLine: '*Get-ItemProperty*'
    SELECTION_8:
        CommandLine:
        - '*HKCU:*'
        - '*HKLM:*'
        - '*registry::*'
        - '*HKEY_*'
    condition: (SELECTION_1 and SELECTION_2 and SELECTION_3 and SELECTION_4 and SELECTION_5
        and SELECTION_6 and SELECTION_7 and SELECTION_8)
falsepositives:
- Unknown
id: c4eeeeae-89f4-43a7-8b48-8d1bdfa66c78
level: high
logsource:
    category: process_creation
    product: windows
modified: 2022/03/15
references:
- https://thedfirreport.com/2022/02/21/qbot-and-zerologon-lead-to-full-domain-compromise/
status: test
tags:
- attack.execution
- attack.persistence
- attack.t1053.005
- attack.t1059.001

